# airootfs/etc/skel/.config/autostart/cinatheme.desktop
[Desktop Entry]
Name=Custom Cinnamon Theme
Exec=cinatheme
Icon=utilities-terminal
Type=Application
Categories=System;Utility;
