# airootfs/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Name=Background Conky
Exec=conky -d
Icon=utilities-terminal
Type=Application
Categories=System;Utility;
